# config/systemd/auto-zram-monitorRepeat.service
[Unit]
Description=Auto-ZRAM monitor

[Service]
Type=oneshot
ExecStart=/usr/bin/auto-zram compressionAnalysisJson

[Install]
WantedBy=timers.target
